# docker/supervisord.conf
[program:nginx]
command=nginx -c /serve/nginx.conf

[program:twisted-web]
command=python /serve/server.py

[supervisord]
childlogdir=/serve
