# docker/supervisord.conf
[program:nginx]
command=nginx -c /data/nginx.conf

[supervisord]

