# docker/supervisord.conf
[program:nginx]
command=nginx -c /serve/nginx.conf

[supervisord]

